# site-cookbooks/manage/templates/default/supervisor/hubot.conf.erb
[program:Hubot]
command=/home/remper/hubot/current/bin/hubot -a lingr
directory=/home/remper/hubot/current
environment=HOME="/home/remper",USER="remper"
user=remper
autorestart=true
stdout_logfile=/var/log/supervisor/hubot.log
redirect_stderr=true

