# site-cookbooks/manage/templates/default/supervisor/hubot.conf.erb
[program:Hubot]
command=cd /home/remper/hubot/current && /bin/hubot -a lingr
user=remper
autorestart=true
stdout_logfile=/var/log/supervisor/hubot.log
redirect_stderr=true

